# proxmox/roles/install_network_checkers_cts/files/supervisor-fast-api.conf
[program:fastapi-app]
command=/opt/fastapi/gunicorn_start
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/null